# deploy/supervisor/conf.d/eval.conf
[program:eval]
#command=bash /home/zeta3/SageMath/sage -python /home/zeta3/numberdb/numberdb-website/services/eval.py
command=sage -python services/eval.py
environment=HOME=/home/numberdb,USER=numberdb,PATH="$PATH:$HOME/bin:/usr/bin:/home/numberdb/SageMath",SAGE_ROOT="/home/numberdb/SageMath"
user=numberdb
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/numberdb/logs/eval.log
directory=/home/numberdb/numberdb-website
